# files/shareVolume/config/supervisor/fluxbox.ini
[program:fluxbox]
environment=LANG=zh_CN.UTF-8
command=fluxbox -display :1 -screen 0
autostart=true
autorestart=true
priority=30
startsecs=5