# files/shareVolume/config/supervisor/jupyter.ini
[program:jupyter]
Environment=LANG=zh_CN.UTF-8
command=jupyter notebook --no-browser --allow-root --ip=0.0.0.0 --port=80
autostart=true
autorestart=true
priority=60