# files/shareVolume/config/supervisor/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf -g 'daemon on;'
autostart=true
autorestart=true
priority=10