# files/shareVolume/config/supervisor/noVNC.ini
[program:noVNC]
command=websockify --web /noVNC/ 80 0.0.0.0:5900
autostart=true
autorestart=true
priority=50