# files/shareVolume/config/supervisor/xray.ini
[program:xray]
command=xray -config /shareVolume/config/xray/config.json
autostart=true
autorestart=true
priority=80